# systemd_copy/helloworld.service
[Unit]
Description="Hello world script"

[Service]
ExecStart=/home/stewells/AfricaNowcasting/dev/systemd/helloworld.sh 
Type=oneshot
